# fedora-silverblue/etc/profile.d/vim-default-editor.sh
# Ensure Vi(m) is set as EDITOR
export EDITOR="/usr/bin/vi"

# Editor to use when editing units;
# overrides $EDITOR and $VISUAL
export SYSTEMD_EDITOR="/usr/bin/vi"
